# etc/cron.d/dotfiles
#!/bin/bash

set -x

REPO=/home/user/dotfiles

cd $REPO

for item in $(find ./ -type f ! -path '*/.git*'); do
  cp -a --parents /$item ${REPO}/
done

